# config/system/pikachu/udev/rules.d/90-set-default-battery-max-charge.rules
ACTION=="add", RUN+="/home/w0ryn/.wryn/bin/pikachu/set-battery-max-charge"
