# config/local/pikachu/etc/udev/rules.d/90-set-default-battery-max-charge.rules
ACTION=="add", RUN+="/home/w0ryn/.local/share/wryn/config/local/pikachu/bin/set-battery-max-charge"
